# deploy/systemd/manshot-cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Manshot Dashboard
After=network.target manshot-dashboard.service
Requires=manshot-dashboard.service

[Service]
Type=simple
User=mandre
Group=mandre
Environment=NO_AUTOUPDATE=1
ExecStart=/usr/local/bin/cloudflared --config /home/mandre/.cloudflared/config.yml tunnel run manshot-dashboard
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
